# external/board/amisr-ccm-v2/overlay/etc/profile.d/ccm.sh
export PATH=$PATH:/opt/transport/groups/ccm/support/bin
export PYTHONPATH=/opt/transport/groups/ccm/support/lib
