# external/board/amisr-ccm-v2/overlay/etc/profile.d/transport.sh
export PATH=$PATH:/opt/transport/bin
export DATA_TRANSPORT_PATH=/opt/transport
