# external/board/amisr-gateway-v1/overlay/etc/profile.d/ccm.sh
export PATH=$PATH:/opt/transport/groups/ccm/support/bin
export PYTHONPATH=/opt/transport/groups/ccm/support/lib
export DATA_TRANSPORT_PATH=/home/ops/demo
